# patches/g_multi.conf
#
# g_multi - g_multi start script
#

author      "Ettore Chimenti"
description "USB Gadget Script"

start on local-filesystems
stop on deconfiguring-networking

pre-start script
  /usr/sbin/g_multi.sh
  test -c /dev/ttyGS0 || { stop; exit 0; }
end script

exec /sbin/getty --noclear 115200 ttyGS0

post-stop script
  for i in `lsmod | grep -o -e \^g_multi -e \^g_serial -e \^g_ether` 
    do rmmod $i
  done
  /etc/init.d/udhcpd stop
end script