# patches/g_multi.conf
#
# g_multi - g_multi start script
#

author      "Ettore Chimenti"
description "USB Gadget Script"

start on local-filesystems
stop on deconfiguring-networking

exec /usr/sbin/g_multi.sh

