# patches/ttymxc0.conf
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn

pre-start script 

  test -c /dev/ttymxc0 || { stop; exit 0; }

end script

exec /sbin/getty --noclear 115200 ttymxc0
